# dist/homekid
#! /bin/sh

### BEGIN INIT INFO
# Provides:		homekid
# Required-Start:	$local_fs $remote_fs postgresql telldusd
# Required-Stop:	$local_fs $remote_fs postgresql telldusd
# Should-Start:		
# Should-Stop:		
# Default-Start:	2 3 4 5	
# Default-Stop:		0 1 6
# Short-Description:	Service for gathering sensors and devices
# Description:		Service for gathering sensors and devices
### END INIT INFO

PATH="/sbin:/usr/sbin:/bin:/usr/bin"
NAME="homeki"
CURRENT_DIR="/opt/homeki"
PIDFILE="/opt/homeki/homeki.pid"
SCRIPTNAME="/etc/init.d/homekid"
USER="homeki"

# Load the VERBOSE setting and other rcS variables
. /lib/init/vars.sh

# Define LSB log_* functions.
# Depend on lsb-base (>= 3.0-6) to ensure that this file is present.
. /lib/lsb/init-functions

#
# Function that starts the daemon/service
#
do_start()
{
	start-stop-daemon --start --background --chdir $CURRENT_DIR --chuid $USER -m -p $PIDFILE --exec /opt/homeki/start.sh \
		|| return 2
}

#
# Function that stops the daemon/service
#
do_stop()
{
	start-stop-daemon --stop --oknodo --pidfile $PIDFILE
	RETVAL="$?"
	rm -f $PIDFILE
	return "$RETVAL"
}

case "$1" in
  start)
	[ "$VERBOSE" != no ] && log_daemon_msg "Starting $NAME" "$NAME"
	do_start
	case "$?" in
		0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
		2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
	esac
	;;
  stop)
	[ "$VERBOSE" != no ] && log_daemon_msg "Stopping $NAME" "$NAME"
	do_stop
	case "$?" in
		0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
		2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
	esac
	;;
  restart)
	#
	# If the "reload" option is implemented then remove the
	# 'force-reload' alias
	#
	log_daemon_msg "Restarting $NAME" "$NAME"
	do_stop
	case "$?" in
	  0|1)
		do_start
		case "$?" in
			0) log_end_msg 0 ;;
			1) log_end_msg 1 ;; # Old process is still running
			*) log_end_msg 1 ;; # Failed to start
		esac
		;;
	  *)
	  	# Failed to stop
		log_end_msg 1
		;;
	esac
	;;
  *)
	echo "Usage: $SCRIPTNAME {start|stop|restart}" >&2
	exit 3
	;;
esac

:
