# dist/postinst
#!/bin/bash

set -e

if [ "$1" = "configure" ] ; then
	echo -n "creating homeki user: "
	if [ "X`grep "^homeki:" /etc/passwd`" == "X" ]; then
		useradd --home-dir /home/homeki --create-home --groups sudo,plugdev homeki
		echo "done"
	else
		echo "not needed, already exists"
	fi
	
	# needed because psql commands rely on fail exit codes to determine what to do
	set +e
	
	echo -n "creating postgresql homeki user: "
	echo "SELECT usename FROM pg_user WHERE usename = 'homeki';" | sudo -u postgres psql | grep homeki > /dev/null 2>&1
	if [ $? -ne 0 ]; then
		echo "CREATE USER homeki WITH PASSWORD 'homeki';" | sudo -u postgres psql > /dev/null 2>&1
		echo "done"
	else
		echo "not needed, already exists"
	fi

	echo -n "creating postgresql homeki database: "
	echo "SELECT datname FROM pg_database WHERE datname = 'homeki';" | sudo -u postgres psql | grep homeki > /dev/null 2>&1
    if [ $? -ne 0 ]; then
        echo "CREATE DATABASE homeki WITH owner = homeki;" | sudo -u postgres psql > /dev/null 2>&1
		echo "done"
	else
		echo "not needed, already exists"
    fi

	set -e
	
	echo -n "changing owner of /mnt/1wire: "
	if [ -d /mnt/1wire ]; then
		chown -R homeki:homeki /mnt/1wire
		echo "done"
	else
		echo "does not exist"
	fi

	echo -n "creating folder /opt/homeki/log: "
	mkdir -p /opt/homeki/log
	echo "done"

	echo -n "changing owner and permissions in /opt/homeki: "
	chown -R homeki:homeki /opt/homeki
	chmod a+w /opt/homeki/stdout.log
	chmod a+w /opt/homeki/stderr.log
	echo "done"
	
	# Remove old left-over S99homeki from rc2.d
	# (can be removed when all servers are updated)
	rm -f /etc/rc2.d/S99homeki

	update-rc.d homekid defaults 99

	update-alternatives --set java /usr/lib/jvm/java-7-openjdk-armhf/jre/bin/java
	
	echo "------------------- APT: STARTED HOMEKI -------------------" >> /opt/homeki/stdout.log
	/etc/init.d/homekid start
fi

exit 0
